# applications/luci-njit/root/etc/uci-defaults/luci-njit
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@njit[-1]
	add ucitrack njit
	set ucitrack.@njit[-1].init=njit
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
